# HXP/PWN/one_byte_unsolved/filesystem/etc/inittab
::sysinit:/etc/init.d/rc.S
::once:-sh -c 'clear ; eval $(resize) ; cat /etc/motd ; setuidgid 0 /bin/sh ; poweroff'
